# mattermost
#!/sbin/runscript
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

name="MatterMost"
description="Mattermost messanger"

mattermost_root=/opt/mattermost
mattermost_group=mattermost
mattermost_user=mattermost
daemon="$mattermost_root/bin/platform"
pidfile=/var/run/mattermost.pid

depend() {
	need net
	use mysql dns logger
	after sshd
}

start() {
    ebegin "Starting mattermost"
    start-stop-daemon --start --quiet \
        --user $mattermost_user --group $mattermost_group --chdir $mattermost_root --background \
        --exec $daemon 
   eend $? 
}

stop() {
    ebegin "Stopping mattermost"
    start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 \
        --exec $daemon
    eend $?
}
